# debian/bewitchd.service
[Unit]
Description=Bewitch Server Monitor Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/bewitchd -config /etc/bewitch.toml
Restart=on-failure
RestartSec=5
User=bewitch
Group=bewitch
SupplementaryGroups=disk video render
StateDirectory=bewitch
RuntimeDirectory=bewitch
AmbientCapabilities=CAP_SYS_RAWIO CAP_PERFMON
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_PERFMON

[Install]
WantedBy=multi-user.target
